# backup/2025-12-10/simple-logs-collect.service
[Unit]
Description=Simple Logs Collection Service

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 %h/simple-logs/log-collector-enhanced.py
WorkingDirectory=%h/simple-logs
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target