# backup/2025-12-10/simple-logs-alerts.timer
[Unit]
Description=Check for alerts every hour

[Timer]
OnCalendar=hourly
Persistent=true
RandomizedDelaySec=600

[Install]
WantedBy=timers.target