# backup/2025-12-10/simple-logs-collect.timer
[Unit]
Description=Collect logs every 30 minutes

[Timer]
OnCalendar=*:0/30
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target