# backup/2025-12-10/simple-logs-git.timer
[Unit]
Description=Upload to GitHub every 2 hours

[Timer]
OnCalendar=0/2:00
Persistent=true

[Install]
WantedBy=timers.target